# meat-price-spreads-viz.service
[Unit]
Description=USDA Meat Price Spreads Dash App
After=network.target

[Service]
User=user
Restart=always
RestartSec=5
ExecStartPre=-/usr/bin/docker stop meat-price-spreads-viz
ExecStartPre=-/usr/bin/docker rm meat-price-spreads-viz
# No API key needed — the app only reads data/, and fetch_data.py pulls public
# ERS file downloads with no credentials
# Served at root on its own Funnel port (8443) — path-based routing was tried
# first but hit an open Tailscale bug (--set-path breaks apps with absolute-
# path sub-resources, see OPERATIONS.md). app.py still supports
# DASH_URL_BASE_PATHNAME if that's ever revisited via a proper reverse proxy.
ExecStart=/usr/bin/docker run --rm --pull=never --name meat-price-spreads-viz \
    -v /home/user/meat-price-spreads-viz/data:/app/data:rw \
    -p 127.0.0.1:8052:8052 \
    localhost/meat-price-spreads-viz:latest
ExecStop=/usr/bin/docker stop meat-price-spreads-viz

[Install]
WantedBy=multi-user.target
