# meat-price-spreads-viz.service
[Unit]
Description=USDA Meat Price Spreads Dash App
After=network.target

[Service]
User=user
Restart=always
RestartSec=5
ExecStartPre=-/usr/bin/docker stop meat-price-spreads-viz
ExecStartPre=-/usr/bin/docker rm meat-price-spreads-viz
# No API key needed — the app only reads data/, and fetch_data.py pulls public
# ERS file downloads with no credentials
ExecStart=/usr/bin/docker run --rm --pull=never --name meat-price-spreads-viz \
    -v /home/user/meat-price-spreads-viz/data:/app/data:rw \
    -p 127.0.0.1:8052:8052 \
    localhost/meat-price-spreads-viz:latest
ExecStop=/usr/bin/docker stop meat-price-spreads-viz

[Install]
WantedBy=multi-user.target
